# devices/betafpv-vr04-goggle/overlay/etc/init.d/ml-display
#!/sbin/openrc-run
# Display bring-up at boot: load the DRM display stack, start the ml-drmfd DRM-master
# broker (the session anchor every display client attaches to: ml-splash now,
# ml-pipeline/ml-hud later), paint the boot splash (the vendor mountain, nosignal.yuv,
# I420 on the DRM primary - covered by video/OSD planes later, uncovered when they
# stop), and switch the backlight on.
#
# Best-effort throughout: missing binaries/asset/drivers log a warning and the boot
# continues - a dark panel then means a real display-driver problem (which is the
# point: splash visible = display chain healthy, before any video work).
#
# `after modules mdev hwdrivers`: the coldplug loads the display chain (VO/DSI/panel/GPIO/PWM), so
# we only wait for /dev/dri/card0, then start the broker + splash + backlight. This service loads no
# modules.
description="Display bring-up: DRM stack + boot splash + backlight"

BROKER=/usr/local/bin/ml-drmfd
SPLASH=/usr/local/bin/ml-splash
IMAGE=/usr/local/share/nosignal.yuv
PIDFILE=/run/ml-drmfd.pid
SOCK=/run/missinglynk/drm.sock
BL=/sys/class/backlight/backlight
HUD_SETTINGS=/usrdata/hud/settings.json

# Set the brightness level from the persisted HUD setting (goggles.brightness, a percent)
# using ml-hud's percent-to-level mapping, so the first light-up is already at the user's
# level and ml-hud's later re-apply is a no-op. No setting (or no file) keeps the driver
# default.
apply_saved_brightness() {
	local pct level max
	pct="$(sed -n 's/.*"brightness"[^0-9]*\([0-9][0-9]*\).*/\1/p' "$HUD_SETTINGS" 2>/dev/null | head -n1)"
	[ -n "$pct" ] || return 0
	max="$(cat "$BL/max_brightness" 2>/dev/null)" || return 0
	level=$((max * pct / 100))
	[ "$level" -ge 1 ] || level=1
	echo "$level" > "$BL/brightness" 2>/dev/null || true
}

depend() {
	after modules mdev hwdrivers
}

start() {
	ebegin "Bringing up display"

	# The coldplug binds artosyn_gpio (the panel reset line) and the display chain
	# from their DT nodes; wait for it to finish.
	local i=0
	while [ ! -e /dev/dri/card0 ] && [ $i -lt 50 ]; do
		sleep 0.1
		i=$((i + 1))
	done

	if [ ! -e /dev/dri/card0 ]; then
		ewarn "ml-display: /dev/dri/card0 never appeared (display drivers missing?)"
		eend 0
		return 0
	fi

	# The kernel never touches the backlight (the panel driver does not attach it), so it is
	# still in its probe-time OFF state here. The explicit off is belt-and-braces for restarts,
	# where a previous run already lit it. The sysfs node exists only once the panel/backlight
	# driver has bound (which card0 above implies).
	echo 1 > "$BL/bl_power" 2>/dev/null

	if [ -x "$BROKER" ]; then
		start-stop-daemon --start --background --make-pidfile \
			--stdout /var/log/ml-drmfd.log --stderr /var/log/ml-drmfd.log \
			--pidfile "$PIDFILE" --exec "$BROKER"
		i=0
		while [ ! -S "$SOCK" ] && [ $i -lt 50 ]; do
			sleep 0.1
			i=$((i + 1))
		done
	else
		ewarn "ml-display: $BROKER missing - no broker, no splash"
	fi

	if [ -x "$SPLASH" ] && [ -S "$SOCK" ]; then
		if [ -f "$IMAGE" ]; then
			"$SPLASH" "$IMAGE" || ewarn "ml-display: splash failed"
		else
			ewarn "ml-display: $IMAGE missing (stage the vendor dump into the rootfs build)"
		fi
	fi

	# First light: splash is committed, set the user's brightness, then enable. Order matters -
	# brightness before bl_power, so the panel never shows the driver-default level.
	apply_saved_brightness
	echo 0 > "$BL/bl_power" 2>/dev/null

	eend 0
}

stop() {
	ebegin "Stopping ml-drmfd (drops the display session)"
	start-stop-daemon --stop --pidfile "$PIDFILE" 2>/dev/null
	eend 0
}
