# devices/betafpv-vr04-goggle/overlay/etc/init.d/ml-video
#!/sbin/openrc-run
# RF video at boot: bring the AR8030 link up (ml-rf-bringup: reset release + SDIO re-probe +
# firmware download + sdio0 config), start ml-linkd (association + READY-gated handshakes), then
# the standalone ml-pipeline via ml-video-up. Best-effort throughout: any missing piece logs a
# warning and the boot continues.
description="RF video: AR8030 bring-up + ml-linkd + ml-pipeline"

RF_BRINGUP=/usr/local/bin/ml-rf-bringup
RF_FW=bb_demo_gnd_d.img                     # baseband firmware in /lib/firmware (staged by rootfs build)
RF_CFG_RACE=bb_config_gnd.json.usr_cfg.json     # race band, 16 channels: its ground-side config, same dir
RF_CFG_NORMAL=bb_config_gnd.json.normal_cfg.json  # normal band, 3 channels: same blob, chan_valid_bmp rewritten
RF_BAND_MARKER=/usrdata/missinglynk/rf-band     # "normal" selects the normal band; absent/anything else = race
LINKD=/usr/local/bin/ml-linkd
VIDEO_UP=/usr/local/bin/ml-video-up

depend() {
	after modules mdev hwdrivers ml-display
}

# The band (Normal/Race) is the config's chan_valid_bmp and only reaches the chip at firmware
# upload, so it is chosen here, once per boot, and a band change needs a reboot to take effect.
rf_cfg() {
	if [ "$(cat "$RF_BAND_MARKER" 2>/dev/null)" = "normal" ] &&
	   [ -f "/lib/firmware/$RF_CFG_NORMAL" ]; then
		echo "$RF_CFG_NORMAL"
	else
		echo "$RF_CFG_RACE"
	fi
}

rf_up() {
	# The whole AR8030 sequence lives in ml-rf-bringup (reset release, SDIO re-probe, firmware
	# download, sdio0 config). It is idempotent - a no-op if artosyn_sdio is already loaded, which
	# is what must happen on a service restart (that driver must NEVER be warm-reloaded: it hangs
	# the device).
	if [ ! -x "$RF_BRINGUP" ]; then
		ewarn "ml-video: $RF_BRINGUP missing - no RF bring-up"
		return 1
	fi

	cfg=$(rf_cfg)
	einfo "ml-video: RF band config $cfg"
	"$RF_BRINGUP" "$RF_FW" "$cfg"
}

start() {
	ebegin "Starting RF video"

	rf_up || ewarn "ml-video: RF bring-up incomplete (video will not flow; SSH + serial still work)"

	if [ -x "$LINKD" ]; then
		start-stop-daemon --start --background --make-pidfile \
			--pidfile /run/ml-linkd.pid --exec "$LINKD" \
			--stdout /var/log/ml-linkd.log --stderr /var/log/ml-linkd.log
	else
		ewarn "ml-video: $LINKD missing"
	fi

	if [ -x "$VIDEO_UP" ]; then
		start-stop-daemon --start --background --make-pidfile \
			--pidfile /run/ml-pipeline.pid \
			--stdout /var/log/ml-pipeline.log --stderr /var/log/ml-pipeline.log \
			--exec "$VIDEO_UP" -- rf
	else
		ewarn "ml-video: $VIDEO_UP missing"
	fi

	eend 0
}

stop() {
	ebegin "Stopping RF video (link/driver stay up - never warm-reload artosyn_sdio)"
	start-stop-daemon --stop --pidfile /run/ml-pipeline.pid 2>/dev/null
	start-stop-daemon --stop --pidfile /run/ml-linkd.pid 2>/dev/null
	eend 0
}
