# skeleton/etc/init.d/ml-sdcard
#!/sbin/openrc-run
# Mount the microSD card at /mnt/sdcard at boot. Ordered after the coldplug (dw_mci-artosyn has bound
# /dev/mmcblk*) and before ml-hud (so the OSD free-space field is populated). Post-boot insert and
# removal are handled by the ml-sdmount mdev hotplug hook, which shares the same mount logic. Boot-time
# only here; best-effort - no card or a failed mount leaves /mnt/sdcard empty and never blocks boot.
description="Mount the microSD card at /mnt/sdcard"

depend() {
	after modules mdev hwdrivers
	before ml-hud
}

start() {
	ebegin "Mounting SD card"
	/usr/local/bin/ml-sdmount
	eend 0
}

stop() {
	ebegin "Unmounting SD card"
	umount -l /mnt/sdcard 2>/dev/null
	eend 0
}
